# crontab
0 0 * * 0 python /bolsometro/db-clear.py